# config/systemd-user/keystrel-daemon.service
[Unit]
Description=GPU faster-whisper Keystrel daemon
After=default.target

[Service]
Type=simple
EnvironmentFile=-%h/.config/keystrel-daemon.env
ExecStart=%h/.local/bin/keystrel-daemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
